# bin/mridb
#!/bin/sh
#
# chkconfig: - 80 20
# description: MRIdb web application

prog="mridb"

. /etc/rc.d/init.d/functions

export JAVA_HOME=/opt/jre1.6.0_30
PLAY=/opt/play-1.2.4/play
APP_PATH=/opt/MRIdb/

run(){
    $PLAY $1 $APP_PATH --%prod --pid_file=/var/lock/subsys/$prog -Dprecompiled=true >/dev/null
    RETVALUE=$? 
    if [ $RETVALUE -eq 0 ]
    then
        echo_success
    else
        echo_failure
    fi
    echo
}

case "$1" in
start)
    echo -n $"Starting $prog: "
    run $1
    ;;  
stop)
    echo -n $"Stopping $prog: "
    run $1
    ;;
restart)
    echo -n $"Restarting $prog: "
    run $1
    ;;
*)
    echo "Usage: $prog {start|stop|restart}" >&2
    exit 1
    ;;
esac

exit $?